# approvebot/service.conf
[program:pmapprovebot]
directory=PUT YOUR BOT'S DIRECTORY HERE
command=/path/to/your/venv/bin/python bot.py AcidMarxism
autorestart=true
user=scott
redirect_stderr=true
